# luci-app-forkop/root/etc/uci-defaults/50_luci-forkop
#!/bin/sh

/usr/bin/forkop luci_postinst >/dev/null 2>&1 || true
exit 0
