# luci-app-tachyon/root/etc/uci-defaults/50_luci-tachyon
#!/bin/sh

/usr/bin/tachyon luci_postinst >/dev/null 2>&1 || true
exit 0
